# MarkDownResources/如何设置jar包后台运行并开机自启动/collect-server.service
[Unit]
Description=collect-server
After=syslog.target network.target

[Service]
Type=simple
ExecStart=/home/jieneng/Docker/hpvm-collect-server/collect-server/build/libs/wdog.sh
ExecStop=/bin/kill -9 $MAINPID

[Install]
WantedBy=multi-user.target
